# update_job.service
[Unit]
Description=System Update
[Service]
Type=simple
ExecStart=/root/bin/script.sh
[Install]
WantedBy=muti-user.target
